# upstart.conf
respawn
respawn limit 15 5
 
start on runlevel [2345]
stop on runlevel [06]
 
script
/usr/bin/python2.7 /var/www/collegecat/app.py 2>&1 >> /var/log/collegecat.log
end script
